# package/DEBIAN/postinst
#!/bin/bash

systemctl enable fan-control
systemctl start fan-control
