# package/DEBIAN/postinst
systemctl enable fan-control
systemctl start fan-control
