# package/DEBIAN/postrm
systemctl stop fan-control.service
systemctl disable fan-control.service
