# package/DEBIAN/prerm
#!/bin/bash

systemctl stop fan-control
systemctl disable fan-control
